# .conf/supervisor.http.conf
[program:goloka-http]
command="/srv/venv/bin/gunicorn" "--keep-alive" "1800" "--enable-stdio-inheritance" "-w" "1" "-b" "ec2-54-218-234-227.us-west-2.compute.amazonaws.com:80" "-k" "socketio.sgunicorn.GeventSocketIOWorker" "goloka.server:app"
environment=PATH="/src/venv/bin:/srv/goloka/bin:$PATH",PORT="80",PYTHONPATH="/srv/goloka:/src/venv:/src/venv/lib/python2.7:$PYTHONPATH",MYSQL_URI="mysql://goloka:vrindavan@golokadb.ck3pwjuvjqri.us-west-2.rds.amazonaws.com:3306/golokadb",REDIS_URI="redis://localhost:6379",LOGLEVEL="DEBUG",DOMAIN="ec2-54-218-234-227.us-west-2.compute.amazonaws.com",HOST="ec2-54-218-234-227.us-west-2.compute.amazonaws.com"
autostart=true
autorestart=true
user=ubuntu
stdout_logfile=/var/log/supervisor/http-stdout.log
stderr_logfile=/var/log/supervisor/http-stderr.log
stdout_logfile_backups=1
stderr_logfile_backups=1
